# tools/orchestrator/com.mira.graphify-refresh.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.mira.graphify-refresh</string>

    <!-- Key-free structural refresh of the graphify code graph against
         origin/main, committed + pushed to main. Fixes the drift: render.py
         (orchestrator-pulse) only pushes HTML, never commits graph.json.
         Script works in an isolated worktree (never touches ~/MIRA's tree)
         and aborts cleanly on any failure. Semantic (LLM) edges are preserved
         from the last full extract; this refreshes AST structure only. -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/charlienode/mira-kg-refresh.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/charlienode/MIRA</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/charlienode/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/charlienode</string>
    </dict>

    <!-- 03:30 local daily — quiet window, after eval-fixer (01:00). -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>3</integer>
        <key>Minute</key>
        <integer>30</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/tmp/mira-kg-refresh.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/mira-kg-refresh.err</string>

    <key>KeepAlive</key>
    <false/>
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
